# udev/rules.d/disc-detection.rules
ACTION=="change", SUBSYSTEMS=="block", KERNEL=="sr[0-9]*", ENV{ID_CDROM}=="?*", RUN+="/opt/docker/marvyn/bin/rip-disc '/dev/%k'"
