# etc/udev/rules.d/disc-detection.rules
ACTION=="change", SUBSYSTEMS=="block", KERNEL=="s[rg][0-9]*", ENV{ID_CDROM}=="?*", RUN+="/opt/docker/marvyn/rip-disc"
